# rtc-ctl/rtc-ctl.cron
#!/bin/sh
#
# Simple cron script - save the current clock periodically in case of
# a power failure or other crash

if (command -v rtc-ctl >/dev/null 2>&1) ; then
  rtc-ctl --quiet --set-clock
fi
